# _SHARED/Scripts/com.brandos.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brandos.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>$HOME/Library/Mobile Documents/com~apple~CloudDocs/CoWork_Projects/_SHARED/Scripts/sync_to_github.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>$HOME/Library/Mobile Documents/com~apple~CloudDocs/CoWork_Projects</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/brandos-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/brandos-sync-error.log</string>
</dict>
</plist>
